# recent-messages2.service
[Unit]
Description=Recent Messages service
Documentation=https://github.com/3R01/recent-messages
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=5min
StartLimitBurst=5

[Service]
Type=simple
User=recent-messages
Group=recent-messages
UMask=0077
WorkingDirectory=/opt/recent-messages/current
Environment=RUST_LOG=twitch_irc=info,recent_messages2=info
ExecStart=/opt/recent-messages/current/bin/recent-messages2 --config /etc/recent-messages/config.toml
Restart=on-failure
RestartSec=5s
TimeoutStopSec=120s
KillSignal=SIGTERM

StateDirectory=recent-messages
StateDirectoryMode=0750
RuntimeDirectory=recent-messages
RuntimeDirectoryMode=0750
LimitNOFILE=524288

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectHostname=true
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictNamespaces=true
ReadWritePaths=/var/lib/recent-messages

[Install]
WantedBy=multi-user.target
